# deployment/upstart/codeshelver.conf
description "Codeshelver Node.js Server"
author      "Dennis Reimann - https://www.codeshelver.com"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production
env DIR=/var/www/codeshelver/current

script
  exec start-stop-daemon --start --chuid deploy:deploy --make-pidfile --pidfile $DIR/log/server.pid --exec /usr/bin/coffee $DIR/server.coffee >> $DIR/log/server.log 2>&1
end script
